# .s6/services/someblocks/run
#!/bin/execlineb
foreground { s6-svwait "$S6_SCANDIR/somebar" }
someblocks
