# .s6/services/someblocks/run
#!/bin/sh
s6-svwait "$S6_SCANDIR/somebar" || exit 1
exec someblocks
