# .s6/services/mqtteer/run
#!/bin/sh
. "$HOME/.mqtteer_env"
exec mqtteer
