# .s6/services/swayidle/run
#!/bin/execlineb -P
importas SWAYIDLE_LOCK SWAYIDLE_LOCK
importas SWAYIDLE_TURN_OFF SWAYIDLE_TURN_OFF

swayidle -w
timeout ${SWAYIDLE_LOCK} "swaylock -f"
timeout ${SWAYIDLE_LOCK} "gpg-connect-agent reloadagent /bye"
timeout ${SWAYIDLE_LOCK} "pkill rbw-agent"
timeout ${SWAYIDLE_TURN_OFF} "toggle_monitors off"
resume "toggle_monitors on"
after-resume "toggle_monitors on"
before-sleep "swaylock -f"
lock "swaylock -f"
idlehint 1200
