# .s6/services/swayidle/run
#!/bin/execlineb -P
swayidle -w
timeout 300 "swaylock -f"
timeout 600 "toggle_monitors off"
after-resume "toggle_monitors on"
before-sleep "swaylock -f"
lock "swaylock -f"
idlehint 1200
