# .s6/services/snooze-minutely/run
#!/bin/execlineb -P
importas HOME HOME
define PARTS_DIR ${HOME}/.snooze/minutely
define TIMEFILE ${HOME}/.snooze/.timefile_minutely
define TIMEWAIT 60

chrt --idle 0 ionice -c 3 snooze -H* -M* -S* -t $TIMEFILE -T $TIMEWAIT
execlineb -Pc
"foreground { run-parts ${PARTS_DIR} } foreground { run-parts ${PARTS_DIR}/local } touch $TIMEFILE"
