# .config/rc/init.d/gammastep
#!/sbin/openrc-run

name="gammastep"
description="set color temperature of display according to time of day"

command="gammastep"
supervisor=supervise-daemon
