# .config/rc/init.d/lxqt-policykit-agent
#!/sbin/openrc-run
XDG_STATE_HOME="${XDG_STATE_HOME:-${HOME}/.local/state}"
LOGFILE="${XDG_STATE_HOME}/rc/${RC_SVCNAME}.log"

name="lxqt-policykit-agent"
description="lxqt polkit authentication agent"

command="lxqt-policykit-agent"
output_log="${LOGFILE}"
error_log="${LOGFILE}"
supervisor=supervise-daemon
