# .config/rc/init.d/lxqt-policykit-agent
#!/sbin/openrc-run
. "${XDG_CONFIG_HOME}/rc/log.sh"

name="lxqt-policykit-agent"
description="lxqt polkit authentication agent"

command="lxqt-policykit-agent"
supervisor=supervise-daemon
