# .config/rc/init.d/mako
#!/sbin/openrc-run

name="mako"
description="notification daemon for wayland"

command="mako"
supervisor=supervise-daemon
