# .config/rc/init.d/mqtteer
#!/sbin/openrc-run
. "${XDG_CONFIG_HOME}/rc/log.sh"

name="mqtteer"
description="system monitoring"

command="mqtteer"
supervisor=supervise-daemon
