# .config/rc/init.d/mqtteer
#!/sbin/openrc-run

name="mqtteer"
description="system monitoring"

command="mqtteer"
supervisor=supervise-daemon
