# .config/rc/init.d/niri
#!/sbin/openrc-run
. "${XDG_CONFIG_HOME}/rc/log.sh"

depend() {
    need dbus
}

command="niri"
command_args="--session"
supervisor="supervise-daemon"
respawn_max=1

start_pre() {
    export GRAPHICAL_TTY=1
    . ~/.config/environment
}
