# .config/rc/init.d/niri
#!/sbin/openrc-run
. "${XDG_CONFIG_HOME}/rc/log.sh"

depend() {
    need dbus
}

command="niri"
command_args="--session"
command_background="yes"
pidfile="${XDG_RUNTIME_DIR}/niri.pid"

start_pre() {
    export GRAPHICAL_TTY=1
    . ~/.config/environment
}
