# .config/rc/init.d/niri-autoselect-portal
#!/sbin/openrc-run
. "${XDG_CONFIG_HOME}/rc/log.sh"

name="niri-autoselect-portal"
description="niri-autoselect-portal daemon"

supervisor="supervise-daemon"
command="/usr/libexec/niri-autoselect-portal"
command_args="-v"
