# .config/rc/init.d/ntfysub
#!/sbin/openrc-run

name="ntfysub"
description="subscribe to ntfy topic"

command="ntfysub"
command_args="${TOPICS}"
supervisor=supervise-daemon
