# .config/rc/init.d/ntfysub
#!/sbin/openrc-run
XDG_STATE_HOME="${XDG_STATE_HOME:-${HOME}/.local/state}"
LOGFILE="${XDG_STATE_HOME}/rc/${RC_SVCNAME}.log"

name="ntfysub"
description="subscribe to ntfy topic"

command="ntfysub"
command_args="${TOPICS}"
output_log="${LOGFILE}"
error_log="${LOGFILE}"
supervisor=supervise-daemon
