# .config/rc/init.d/otd
#!/sbin/openrc-run
XDG_STATE_HOME="${XDG_STATE_HOME:-${HOME}/.local/state}"
LOGFILE="${XDG_STATE_HOME}/rc/${RC_SVCNAME}.log"

name="otd"
description="Open Tablet Daemon"

command="otd-daemon"
output_log="${LOGFILE}"
error_log="${LOGFILE}"
supervisor=supervise-daemon
