# .config/rc/init.d/otd
#!/sbin/openrc-run
. "${XDG_CONFIG_HOME}/rc/log.sh"

name="otd"
description="Open Tablet Daemon"

command="otd-daemon"
supervisor=supervise-daemon
