# .config/rc/init.d/otd
#!/sbin/openrc-run

name="otd"
description="Open Tablet Daemon"

command="otd-daemon"
#command_args=""
supervisor=supervise-daemon
