# .config/rc/init.d/pipewire-filter-chain
#!/sbin/openrc-run
. "${XDG_CONFIG_HOME}/rc/log.sh"

depend() {
    need pipewire
}

name="pipewire-filter-chain"
description="pipewire filter chain"

command="pipewire"
command_args="-c filter-chain.conf"
supervisor=supervise-daemon
