# .config/rc/init.d/pipewire-filter-chain
#!/sbin/openrc-run
XDG_STATE_HOME="${XDG_STATE_HOME:-${HOME}/.local/state}"
LOGFILE="${XDG_STATE_HOME}/rc/${RC_SVCNAME}.log"

name="pipewire-filter-chain"
description="pipewire filter chain"

command="pipewire"
command_args="-c filter-chain.conf"
output_log="${LOGFILE}"
error_log="${LOGFILE}"
supervisor=supervise-daemon

depend() {
    need pipewire
}
