# .config/rc/init.d/pipewire
#!/sbin/openrc-run
. "${XDG_CONFIG_HOME}/rc/log.sh"

depend() {
    need dbus
}

name="pipewire"
description="pipewire media server"

command="/usr/bin/pipewire"
supervisor=supervise-daemon
