# .config/rc/init.d/pipewire
#!/sbin/openrc-run
LOGFILE="${XDG_STATE_HOME}/rc/${RC_SVCNAME}.log"

name="pipewire"
description="pipewire media server"

command="/usr/bin/pipewire"
output_log="${LOGFILE}"
error_log="${LOGFILE}"
supervisor=supervise-daemon
