# .config/rc/init.d/pipewire-pulse
#!/sbin/openrc-run
# Copyright 1999-2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License, v2 or later

depend() {
	need pipewire
}

supervisor=supervise-daemon
command="/usr/bin/pipewire"
command_args="-c pipewire-pulse.conf"
