# .config/rc/init.d/snooze
#!/sbin/openrc-run
XDG_STATE_HOME="${XDG_STATE_HOME:-${HOME}/.local/state}"
LOGFILE="${XDG_STATE_HOME}/rc/${RC_SVCNAME}.log"

TIMEDEF=${RC_SVCNAME#*.}
SNOOZE_DIR="${HOME}/.snooze"

SNOOZE_RUN_DIR="${SNOOZE_DIR}/${TIMEDEF}"
TIMEFILE="${SNOOZE_DIR}/.timefile_${TIMEDEF}"

case ${TIMEDEF} in
    minutely)
        TIMEWAIT=60
        ;;
    hourly)
        TIMEWAIT=1h
        ;;
    daily)
        TIMEWAIT=1d
        ;;
esac

name="snooze"
description="periodic tasks"

command="snooze"
command_args="-H* -M* -S* -t ${TIMEFILE} -T ${TIMEWAIT} snooze-run ${SNOOZE_RUN_DIR} ${TIMEFILE}"
output_log="${LOGFILE}"
error_log="${LOGFILE}"
supervisor=supervise-daemon
supervise_daemon_args="--ionice 3 --nicelevel 10"
