# .config/rc/init.d/swayidle
#!/sbin/openrc-run

name="swayidle"
description="Idle manager for Wayland"

command="swayidle"
command_args="-w"
supervisor=supervise-daemon

start_pre() {
    export NIRI_SOCKET
}
