# .config/rc/init.d/wireplumber
#!/sbin/openrc-run
# Copyright 1999-2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License, v2 or later
LOGFILE="${XDG_STATE_HOME}/rc/${RC_SVCNAME}.log"

depend() {
	need pipewire dbus
}

supervisor=supervise-daemon
command="/usr/bin/wireplumber"
output_log="${LOGFILE}"
error_log="${LOGFILE}"
