# .config/rc/init.d/wireplumber
#!/sbin/openrc-run
# Copyright 1999-2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License, v2 or later

depend() {
	need pipewire dbus
}

supervisor=supervise-daemon
command="/usr/bin/wireplumber"
