# .config/rc/init.d/wlsunset
#!/sbin/openrc-run
LOGFILE="${XDG_STATE_HOME}/rc/${RC_SVCNAME}.log"
WLSUNSET_OPTS="${WLSUNSET_OPTS:--S 7:00 -s 20:00 -d 60}"

name="wlsunset"
description="set color temperature of display according to time of day"

command="wlsunset"
command_args="${WLSUNSET_OPTS}"
output_log="${LOGFILE}"
error_log="${LOGFILE}"
supervisor=supervise-daemon
