# .config/rc/init.d/xdg-desktop-portal-gnome
#!/sbin/openrc-run
. "${XDG_CONFIG_HOME}/rc/log.sh"

name="xdg-desktop-portal-gnome"
description="xdg-desktop-portal-gnome daemon"

supervisor="supervise-daemon"
command="/usr/libexec/xdg-desktop-portal-gnome"
command_args="-v"
