# .config/rc/init.d/xdg-desktop-portal
#!/sbin/openrc-run
LOGFILE="${XDG_STATE_HOME}/rc/${RC_SVCNAME}.log"

name="xdg-desktop-portal"
description="xdg-desktop-portal daemon"

supervisor="supervise-daemon"
command="/usr/libexec/xdg-desktop-portal"
output_log="${LOGFILE}"
error_log="${LOGFILE}"
