# .config/rc/init.d/xdg-desktop-portal-gtk
#!/sbin/openrc-run
. "${XDG_CONFIG_HOME}/rc/log.sh"

name="xdg-desktop-portal-gtk"
description="xdg-desktop-portal-gtk daemon"

supervisor="supervise-daemon"
command="/usr/libexec/xdg-desktop-portal-gtk"
command_args="-v"
