# .config/rc/init.d/xdg-desktop-portal-lxqt
#!/sbin/openrc-run
. "${XDG_CONFIG_HOME}/rc/log.sh"

name="xdg-desktop-portal-lxqt"
description="xdg-desktop-portal-lxqt daemon"

supervisor="supervise-daemon"
command="/usr/libexec/xdg-desktop-portal-lxqt"
command_args="-v"
