# .config/rc/init.d/xdg-desktop-portal
#!/sbin/openrc-run
. "${XDG_CONFIG_HOME}/rc/log.sh"

name="xdg-desktop-portal"
description="xdg-desktop-portal daemon"

supervisor="supervise-daemon"
command="/usr/libexec/xdg-desktop-portal"
command_args="-v"
