# .config/rc/init.d/xdg-desktop-portal-lxqt
#!/sbin/openrc-run
LOGFILE="${XDG_STATE_HOME}/rc/${RC_SVCNAME}.log"

name="xdg-desktop-portal-lxqt"
description="xdg-desktop-portal-lxqt daemon"

supervisor="supervise-daemon"
command="/usr/libexec/xdg-desktop-portal-lxqt"
output_log="${LOGFILE}"
error_log="${LOGFILE}"
