# .s6/services/pipewire/run
#!/bin/execlineb -P
s6-notifyoncheck gentoo-pipewire-launcher restart
