# .s6/services/pipewire/run
#!/bin/execlineb
s6-notifyoncheck gentoo-pipewire-launcher restart
