# .s6/services/transmission-daemon/run
#!/bin/execlineb -P
transmission-daemon -f
