# .s6/services/transmission-daemon/run
#!/bin/sh
exec chrt --idle 0 ionice -c 3 transmission-daemon -f
