# .s6/services/snooze-daily/run
#!/bin/execlineb -P
importas HOME HOME
define PARTS_DIR ${HOME}/.snooze/daily
define TIMEFILE ${HOME}/.snooze/.timefile_daily
define TIMEWAIT 1d

snooze -H* -M* -S* -t $TIMEFILE -T $TIMEWAIT
execlineb -Pc "foreground { run-parts ${PARTS_DIR} } touch $TIMEFILE"
