# .s6/services/yambar/run
#!/bin/sh
s6-svwait -U ${S6_SCANDIR}/pipewire
sleep 1
exec yambar
