# .s6/services/swaybg/run
#!/bin/execlineb -P
importas HOME HOME
swaybg -i ${HOME}/.wallpaper
