# .s6/services/somebar/run
#!/bin/execlineb -P
importas DWL_STATUS_PIPE DWL_STATUS_PIPE
redirfd -r 0 ${DWL_STATUS_PIPE} somebar
