# .s6/services/somebar/run
#!/bin/execlineb -P
importas XDG_RUNTIME_DIR XDG_RUNTIME_DIR
redirfd -r 0 ${XDG_RUNTIME_DIR}/dwl_status somebar
