# .s6/services/ssh-agent/run
#!/bin/execlineb -P
importas XDG_RUNTIME_DIR XDG_RUNTIME_DIR
ssh-agent -a ${XDG_RUNTIME_DIR}/ssh-agent.default -t 7200 -D
