# .s6/services/otd/finish
#!/bin/execlineb -P
# some processes seem to stay under PID 1
pkill -f OpenTablet
