# .s6/services/foot-server/run
#!/bin/sh
exec foot --server="${XDG_RUNTIME_DIR}/foot.sock"
